# dock/fdsfgs/images/rtorrent-tmux-ssh/etc/service/rtorrent-tmux/finish
#!/bin/sh
su-exec rtorrent /usr/bin/tmux send-keys -t rt:rtorrent.0 C-q; while pidof rtorrent > /dev/null; do echo stopping rtorrent...; sleep 1; done
