# dock/fdsfgs/images/rtorrent/etc/service/rtorrent/run
#!/bin/sh

set -x

# run
exec 2>&1
exec /bin/busybox su -l -c "TERM=xterm #@TODO: ssh command that launch tmux with rtorrent /usr/bin/rtorrent" rtorrent #@TODO: /dev/stdout or nothing?
